# simple-media-sync/media-sync.service
[Unit]
Description=Media Sync Workflow
After=network.target syncthing.service
Wants=syncthing.service

[Service]
Type=simple
User=root
WorkingDirectory=/app
Environment=PATH=/opt/media-sync-env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/media-sync-env/bin/python3 /app/headless_runner.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/app

[Install]
WantedBy=multi-user.target
